# src/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=mbsync -Va
ExecStartPost=notmuch new
RuntimeMaxSec=600
